# centos/v2ray
#!/bin/sh
# chkconfig: 2345 90 10
# description: Starts/Stops the v2ray Server

/etc/v2ray/v2ray
